# GitVerse Tasks/GitVerse_Task_3/test_manual/services/test-service.service
[Service]
ExecStart=python -c "import time; print('Service started'); time.sleep(60)"
Restart=always
